# crontab
13 17 * * wed sh Edit/RSSAC002-data/cron.sh
